# pkg-build/BUILDROOT/Yearning-3.1.7-1.1.el7.x86_64/etc/profile.d/Yearning.sh
#!/bin/bash

YEARNING=/usr/local/Yearning/
PATH=$YEARNING/bin:$PATH

export YEARNING
export PATH
